# services/avahi/run
#!/bin/bash

# Exit with an error until dbus is available
sv start dbus || exit 1

exec avahi-daemon --no-chroot
