# services/airsonos/run
#!/bin/bash

exec /sbin/setuser nobody airsonos
